# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=influxd run -config /etc/influxdb.conf
autorestart=true
autostart=true

[program:telegraf]
command=telegraf
autorestart=true
autostart=true