# srcs/mysql/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=sh start_mysql.sh
autorestart=true
autostart=true

[program:telegraf]
command=telegraf
autorestart=true
autostart=true