# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command=grafana-server
autorestart=true
autostart=true

[program:telegraf]
command=telegraf
autorestart=true
autostart=true